# CTI_reports/report_728_20251223_105012_ground_truth.json
{
    "report_id": "report_728",
    "generated_at": "20251223_105012",
    "scenario_target_ids": [
        "T1597.002",
        "T1595.001",
        "T1591.001",
        "T1195",
        "T1053",
        "T1547.010",
        "T1548.006",
        "T1134.003",
        "T1613",
        "T1119",
        "T1056",
        "T1213.002",
        "T1071.004",
        "T1052.001"
    ],
    "scenario_flow": [
        "reconnaissance",
        "reconnaissance",
        "reconnaissance",
        "initial-access",
        "execution",
        "persistence",
        "privilege-escalation",
        "privilege-escalation",
        "discovery",
        "collection",
        "collection",
        "collection",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1597.002",
        "T1595.001",
        "T1591.001",
        "T1195",
        "T1053",
        "T1547.010",
        "T1548.006",
        "T1134.003",
        "T1613",
        "T1119",
        "T1056",
        "T1213.002",
        "T1071.004",
        "T1052.001"
    ],
    "actual_found_ids": [
        "T1597.002",
        "T1591.001",
        "T1613",
        "T1547.010",
        "T1595.001",
        "T1195",
        "T1548.006",
        "T1134.003",
        "T1056",
        "T1119",
        "T1053",
        "T1071.004",
        "T1213.002",
        "T1052.001"
    ],
    "clean_text": "ShadowGather actors initiated their campaign by purchasing a list of vulnerable systems from a dark web marketplace, which they then used to conduct reconnaissance by scanning IP blocks to identify active network nodes and determine physical locations of targeted executives via social media queries. Subsequently, the actors compromised the software supply chain by altering a continuous integration pipeline to bundle spyware with legitimate updates, enabling initial access. Once inside, they established persistence by creating a scheduled task that executed a remote access tool every 15 minutes and registered a malicious DLL as a port monitor using the AddMonitor API. To escalate privileges, the malware edited the TCC.db file to authorize accessibility access for a keylogger and impersonated a domain admin user token by invoking SetThreadToken. The actors then conducted discovery operations, querying the Docker API to uncover container resource limits and requests while automatically collecting data from compromised systems. For credential theft, they deployed a fake corporate intranet portal that transparently logged all entered credentials and queried SharePoint for password-protected documents containing proprietary information. Exfiltration was achieved through DNS-based command and control communications using IPv4 A and IPv6 AAAA records, and data was moved from airgapped networks to Internet-connected systems via removable USB devices. Throughout the campaign, system logs showed routine network traffic spikes and failed login attempts from unrecognized IP addresses, though these were initially dismissed as false positives. The targeted organization, a fictional energy company with operations across North America, reported no significant disruptions, but post-incident analysis revealed extensive data loss.",
    "tagged_text": "ShadowGather actors initiated their campaign by purchasing a list of vulnerable systems from a dark web marketplace [T1597.002], which they then used to conduct reconnaissance by scanning IP blocks [T1595.001] to identify active network nodes and determine physical locations [T1591.001] of targeted executives via social media queries. Subsequently, the actors compromised the software supply chain by altering a continuous integration pipeline to bundle spyware with legitimate updates [T1195], enabling initial access. Once inside, they established persistence by creating a scheduled task [T1053] that executed a remote access tool every 15 minutes and registered a malicious DLL as a port monitor [T1547.010] using the AddMonitor API. To escalate privileges, the malware edited the TCC.db file to authorize accessibility access for a keylogger [T1548.006] and impersonated a domain admin user token [T1134.003] by invoking SetThreadToken. The actors then conducted discovery operations, querying the Docker API to uncover container resource limits and requests [T1613] while automatically collecting data [T1119] from compromised systems. For credential theft, they deployed a fake corporate intranet portal that transparently logged all entered credentials [T1056] and queried SharePoint for password-protected documents containing proprietary information [T1213.002]. Exfiltration was achieved through DNS-based command and control communications [T1071.004] using IPv4 A and IPv6 AAAA records, and data was moved from airgapped networks to Internet-connected systems via removable USB devices [T1052.001]. Throughout the campaign, system logs showed routine network traffic spikes and failed login attempts from unrecognized IP addresses, though these were initially dismissed as false positives. The targeted organization, a fictional energy company with operations across North America, reported no significant disruptions, but post-incident analysis revealed extensive data loss.",
    "sentence_annotations": [
        {
            "sentence": "Subsequently, the actors compromised the software supply chain by altering a continuous integration pipeline to bundle spyware with legitimate updates, enabling initial access.",
            "technique_id": "T1195",
            "technique_name": "Supply Chain Compromise"
        },
        {
            "sentence": "Once inside, they established persistence by creating a scheduled task that executed a remote access tool every 15 minutes and registered a malicious DLL as a port monitor [T1547.",
            "technique_id": "T1053",
            "technique_name": "Scheduled Task/Job"
        },
        {
            "sentence": "The actors then conducted discovery operations, querying the Docker API to uncover container resource limits and requests while automatically collecting data from compromised systems.",
            "technique_id": "T1613",
            "technique_name": "Container and Resource Discovery"
        },
        {
            "sentence": "The actors then conducted discovery operations, querying the Docker API to uncover container resource limits and requests while automatically collecting data from compromised systems.",
            "technique_id": "T1119",
            "technique_name": "Automated Collection"
        },
        {
            "sentence": "For credential theft, they deployed a fake corporate intranet portal that transparently logged all entered credentials and queried SharePoint for password-protected documents containing proprietary information [T1213.",
            "technique_id": "T1056",
            "technique_name": "Input Capture"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'ShadowGather' Targeting Critical Infrastructure",
        "summary": "This advisory details a simulated campaign by a fictional advanced persistent threat (APT) group, designated 'ShadowGather', targeting critical infrastructure sectors. Based on intelligence from CISA, FBI, and NSA, the actors employ a multi-stage approach involving reconnaissance, supply chain compromise, persistence mechanisms, and data exfiltration to infiltrate networks and steal sensitive information. The report is intended for blue team training to enhance detection and response capabilities against sophisticated attack patterns.",
        "full_report_text": "ShadowGather actors initiated their campaign by purchasing a list of vulnerable systems from a dark web marketplace [T1597.002], which they then used to conduct reconnaissance by scanning IP blocks [T1595.001] to identify active network nodes and determine physical locations [T1591.001] of targeted executives via social media queries. Subsequently, the actors compromised the software supply chain by altering a continuous integration pipeline to bundle spyware with legitimate updates [T1195], enabling initial access. Once inside, they established persistence by creating a scheduled task [T1053] that executed a remote access tool every 15 minutes and registered a malicious DLL as a port monitor [T1547.010] using the AddMonitor API. To escalate privileges, the malware edited the TCC.db file to authorize accessibility access for a keylogger [T1548.006] and impersonated a domain admin user token [T1134.003] by invoking SetThreadToken. The actors then conducted discovery operations, querying the Docker API to uncover container resource limits and requests [T1613] while automatically collecting data [T1119] from compromised systems. For credential theft, they deployed a fake corporate intranet portal that transparently logged all entered credentials [T1056] and queried SharePoint for password-protected documents containing proprietary information [T1213.002]. Exfiltration was achieved through DNS-based command and control communications [T1071.004] using IPv4 A and IPv6 AAAA records, and data was moved from airgapped networks to Internet-connected systems via removable USB devices [T1052.001]. Throughout the campaign, system logs showed routine network traffic spikes and failed login attempts from unrecognized IP addresses, though these were initially dismissed as false positives. The targeted organization, a fictional energy company with operations across North America, reported no significant disruptions, but post-incident analysis revealed extensive data loss.",
        "techniques_used": [
            "T1597.002",
            "T1595.001",
            "T1591.001",
            "T1195",
            "T1053",
            "T1547.010",
            "T1548.006",
            "T1134.003",
            "T1613",
            "T1119",
            "T1056",
            "T1213.002",
            "T1071.004",
            "T1052.001"
        ]
    }
}