# CTI_reports/report_122_20251222_133459_ground_truth.json
{
    "report_id": "report_122",
    "generated_at": "20251222_133459",
    "scenario_target_ids": [
        "T1590.006",
        "T1596.004",
        "T1200",
        "T1053.005",
        "T1552.006",
        "T1550.002",
        "T1001.002",
        "T1048.002",
        "T1030",
        "T1561.001",
        "T1491"
    ],
    "scenario_flow": [
        "reconnaissance",
        "reconnaissance",
        "initial-access",
        "execution",
        "credential-access",
        "lateral-movement",
        "command-and-control",
        "exfiltration",
        "exfiltration",
        "impact",
        "impact"
    ],
    "llm_declared_ids": [
        "T1590.006",
        "T1596.004",
        "T1200",
        "T1053.005",
        "T1552.006",
        "T1550.002",
        "T1001.002",
        "T1048.002",
        "T1030",
        "T1561.001",
        "T1491"
    ],
    "actual_found_ids": [
        "T1491",
        "T1030",
        "T1550.002",
        "T1001.002",
        "T1200",
        "T1590.006",
        "T1561.001",
        "T1596.004",
        "T1048.002",
        "T1053.005",
        "T1552.006"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to report on a simulated APT campaign targeting critical infrastructure sectors, based on a fictional scenario for training. The actors likely conduct reconnaissance to identify network security appliances and use nslookup commands to resolve target domains and identify associated CDN IP ranges, gathering intelligence on network defenses. On 2024-10-01, network administrators at a fictional energy company discovered a passive network tap connected to the main switch in the server room, indicating physical hardware additions. System logs from the server showed unusual scheduled tasks created via PowerShell, with Kapeka malware persisting through scheduled tasks while simultaneously executing reconnaissance scripts that used 'dir /s *.xml' to identify all XML files potentially containing Group Policy Preferences credentials. The company's background includes operations across multiple states, with no prior incidents reported. Following initial access, the adversary employed dumped hashes to authenticate to other machines via pass the hash and used steganography tool.exe to conceal command and control traffic within legitimate-looking images, evading detection. Data exfiltration involved employing OpenSSL to encrypt financial data prior to exfiltration over asymmetric encrypted non-C2 protocols, with compressed log files stored incrementally at https://c2.com/chunk with delays between each upload to manage data transfer size limits. Post-exfiltration, the wiper malware wrote random data patterns across all sectors of the hard drive and the virus injected fake blue screen of death images to panic users into paying ransoms as part of a defacement strategy. Incident response teams noted increased network latency but attributed it to routine maintenance.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to report on a simulated APT campaign targeting critical infrastructure sectors, based on a fictional scenario for training. The actors likely conduct reconnaissance to identify network security appliances [T1590.006] and use nslookup commands to resolve target domains and identify associated CDN IP ranges [T1596.004], gathering intelligence on network defenses. On 2024-10-01, network administrators at a fictional energy company discovered a passive network tap connected to the main switch in the server room, indicating physical hardware additions [T1200]. System logs from the server showed unusual scheduled tasks created via PowerShell, with Kapeka malware persisting through scheduled tasks [T1053.005] while simultaneously executing reconnaissance scripts that used 'dir /s *.xml' to identify all XML files potentially containing Group Policy Preferences credentials [T1552.006]. The company's background includes operations across multiple states, with no prior incidents reported. Following initial access, the adversary employed dumped hashes to authenticate to other machines via pass the hash [T1550.002] and used steganography tool.exe to conceal command and control traffic within legitimate-looking images [T1001.002], evading detection. Data exfiltration involved employing OpenSSL to encrypt financial data prior to exfiltration over asymmetric encrypted non-C2 protocols [T1048.002], with compressed log files stored incrementally at https://c2.com/chunk with delays between each upload to manage data transfer size limits [T1030]. Post-exfiltration, the wiper malware wrote random data patterns across all sectors of the hard drive [T1561.001] and the virus injected fake blue screen of death images to panic users into paying ransoms as part of a defacement strategy [T1491]. Incident response teams noted increased network latency but attributed it to routine maintenance.",
    "sentence_annotations": [
        {
            "sentence": "On 2024-10-01, network administrators at a fictional energy company discovered a passive network tap connected to the main switch in the server room, indicating physical hardware additions.",
            "technique_id": "T1200",
            "technique_name": "Hardware Additions"
        },
        {
            "sentence": "com/chunk with delays between each upload to manage data transfer size limits.",
            "technique_id": "T1030",
            "technique_name": "Data Transfer Size Limits"
        },
        {
            "sentence": "001] and the virus injected fake blue screen of death images to panic users into paying ransoms as part of a defacement strategy.",
            "technique_id": "T1491",
            "technique_name": "Defacement"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Simulated APT Campaign Targeting Critical Infrastructure",
        "summary": "This advisory details a fictional campaign by a simulated advanced persistent threat (APT) actor, designed for blue team training purposes. The campaign involves reconnaissance, initial access via hardware compromise, persistence mechanisms, credential theft, data exfiltration, and destructive actions. Techniques are embedded with MITRE ATT&CK IDs for educational reference.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to report on a simulated APT campaign targeting critical infrastructure sectors, based on a fictional scenario for training. The actors likely conduct reconnaissance to identify network security appliances [T1590.006] and use nslookup commands to resolve target domains and identify associated CDN IP ranges [T1596.004], gathering intelligence on network defenses. On 2024-10-01, network administrators at a fictional energy company discovered a passive network tap connected to the main switch in the server room, indicating physical hardware additions [T1200]. System logs from the server showed unusual scheduled tasks created via PowerShell, with Kapeka malware persisting through scheduled tasks [T1053.005] while simultaneously executing reconnaissance scripts that used 'dir /s *.xml' to identify all XML files potentially containing Group Policy Preferences credentials [T1552.006]. The company's background includes operations across multiple states, with no prior incidents reported. Following initial access, the adversary employed dumped hashes to authenticate to other machines via pass the hash [T1550.002] and used steganography tool.exe to conceal command and control traffic within legitimate-looking images [T1001.002], evading detection. Data exfiltration involved employing OpenSSL to encrypt financial data prior to exfiltration over asymmetric encrypted non-C2 protocols [T1048.002], with compressed log files stored incrementally at https://c2.com/chunk with delays between each upload to manage data transfer size limits [T1030]. Post-exfiltration, the wiper malware wrote random data patterns across all sectors of the hard drive [T1561.001] and the virus injected fake blue screen of death images to panic users into paying ransoms as part of a defacement strategy [T1491]. Incident response teams noted increased network latency but attributed it to routine maintenance.",
        "techniques_used": [
            "T1590.006",
            "T1596.004",
            "T1200",
            "T1053.005",
            "T1552.006",
            "T1550.002",
            "T1001.002",
            "T1048.002",
            "T1030",
            "T1561.001",
            "T1491"
        ]
    }
}