# CTI_reports/report_375_20251222_152039_ground_truth.json
{
    "report_id": "report_375",
    "generated_at": "20251222_152039",
    "scenario_target_ids": [
        "T1593.001",
        "T1592.002",
        "T1593",
        "T1669",
        "T1203",
        "T1546",
        "T1127.002",
        "T1647",
        "T1606.001",
        "T1557.003",
        "T1074",
        "T1074.001",
        "T1102.001",
        "T1567.003"
    ],
    "scenario_flow": [
        "reconnaissance",
        "reconnaissance",
        "reconnaissance",
        "initial-access",
        "execution",
        "privilege-escalation",
        "defense-evasion",
        "defense-evasion",
        "credential-access",
        "credential-access",
        "collection",
        "collection",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1593.001",
        "T1592.002",
        "T1593",
        "T1669",
        "T1203",
        "T1546",
        "T1127.002",
        "T1647",
        "T1606.001",
        "T1557.003",
        "T1074",
        "T1074.001",
        "T1102.001",
        "T1567.003"
    ],
    "actual_found_ids": [
        "T1593",
        "T1593.001",
        "T1127.002",
        "T1102.001",
        "T1606.001",
        "T1546",
        "T1669",
        "T1074.001",
        "T1647",
        "T1567.003",
        "T1074",
        "T1557.003",
        "T1592.002",
        "T1203"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide training insights into a fictional APT campaign. Crimson Viper actors likely conduct reconnaissance by monitoring social media platforms and querying business directories to gather organizational details, while also scanning FTP server banners to identify software versions. During the initial access phase, the actors established wireless connections to enterprise Wi-Fi networks, exploiting vulnerabilities like CVE-2022-30190 for client execution and using rundll32.exe to execute shellcode from compromised DLL files via ClickOnce. For persistence, they modified Plist files to run scripts upon system events and created scheduled tasks to execute malicious DLLs, while maintaining access through forged web cookies. In the network phase, the malware abused DHCP protocol options to spoof DNS queries, redirecting traffic to malicious servers. Data collection involved staging files locally using cmd.exe to save metadata to temporary directories and storing system files in working folders, with payloads accessed via dead drop resolvers. Finally, exfiltration was conducted by uploading collected email archives to text storage sites using PowerShell scripts. System logs from the targeted environment showed routine network traffic and user login attempts, but no anomalies were flagged initially, highlighting the need for enhanced monitoring of these TTPs.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide training insights into a fictional APT campaign. Crimson Viper actors likely conduct reconnaissance by monitoring social media platforms [T1593.001] and querying business directories [T1593] to gather organizational details, while also scanning FTP server banners to identify software versions [T1592.002]. During the initial access phase, the actors established wireless connections to enterprise Wi-Fi networks [T1669], exploiting vulnerabilities like CVE-2022-30190 for client execution [T1203] and using rundll32.exe to execute shellcode from compromised DLL files via ClickOnce [T1127.002]. For persistence, they modified Plist files to run scripts upon system events [T1647] and created scheduled tasks to execute malicious DLLs [T1546], while maintaining access through forged web cookies [T1606.001]. In the network phase, the malware abused DHCP protocol options to spoof DNS queries [T1557.003], redirecting traffic to malicious servers. Data collection involved staging files locally using cmd.exe to save metadata to temporary directories [T1074] and storing system files in working folders [T1074.001], with payloads accessed via dead drop resolvers [T1102.001]. Finally, exfiltration was conducted by uploading collected email archives to text storage sites using PowerShell scripts [T1567.003]. System logs from the targeted environment showed routine network traffic and user login attempts, but no anomalies were flagged initially, highlighting the need for enhanced monitoring of these TTPs.",
    "sentence_annotations": [
        {
            "sentence": "001] and querying business directories to gather organizational details, while also scanning FTP server banners to identify software versions [T1592.",
            "technique_id": "T1593",
            "technique_name": "Search Open Websites/Domains"
        },
        {
            "sentence": "During the initial access phase, the actors established wireless connections to enterprise Wi-Fi networks, exploiting vulnerabilities like CVE-2022-30190 for client execution and using rundll32.",
            "technique_id": "T1669",
            "technique_name": "Wi-Fi Networks"
        },
        {
            "sentence": "During the initial access phase, the actors established wireless connections to enterprise Wi-Fi networks, exploiting vulnerabilities like CVE-2022-30190 for client execution and using rundll32.",
            "technique_id": "T1203",
            "technique_name": "Exploitation for Client Execution"
        },
        {
            "sentence": "For persistence, they modified Plist files to run scripts upon system events and created scheduled tasks to execute malicious DLLs, while maintaining access through forged web cookies [T1606.",
            "technique_id": "T1647",
            "technique_name": "Plist File Modification"
        },
        {
            "sentence": "For persistence, they modified Plist files to run scripts upon system events and created scheduled tasks to execute malicious DLLs, while maintaining access through forged web cookies [T1606.",
            "technique_id": "T1546",
            "technique_name": "Event Triggered Execution"
        },
        {
            "sentence": "exe to save metadata to temporary directories and storing system files in working folders [T1074.",
            "technique_id": "T1074",
            "technique_name": "Data Staged"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign Targeting Financial Sector",
        "summary": "This advisory details a simulated campaign by a fictional APT group, designated as 'Crimson Viper', targeting financial institutions for data exfiltration and network compromise. Based on observed tactics, techniques, and procedures (TTPs), Crimson Viper employs a multi-faceted approach involving reconnaissance, initial access via Wi-Fi exploitation, client-side execution, persistence mechanisms, and data staging for exfiltration. The report is intended for blue team training to enhance detection and response capabilities in real-world scenarios.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide training insights into a fictional APT campaign. Crimson Viper actors likely conduct reconnaissance by monitoring social media platforms [T1593.001] and querying business directories [T1593] to gather organizational details, while also scanning FTP server banners to identify software versions [T1592.002]. During the initial access phase, the actors established wireless connections to enterprise Wi-Fi networks [T1669], exploiting vulnerabilities like CVE-2022-30190 for client execution [T1203] and using rundll32.exe to execute shellcode from compromised DLL files via ClickOnce [T1127.002]. For persistence, they modified Plist files to run scripts upon system events [T1647] and created scheduled tasks to execute malicious DLLs [T1546], while maintaining access through forged web cookies [T1606.001]. In the network phase, the malware abused DHCP protocol options to spoof DNS queries [T1557.003], redirecting traffic to malicious servers. Data collection involved staging files locally using cmd.exe to save metadata to temporary directories [T1074] and storing system files in working folders [T1074.001], with payloads accessed via dead drop resolvers [T1102.001]. Finally, exfiltration was conducted by uploading collected email archives to text storage sites using PowerShell scripts [T1567.003]. System logs from the targeted environment showed routine network traffic and user login attempts, but no anomalies were flagged initially, highlighting the need for enhanced monitoring of these TTPs.",
        "techniques_used": [
            "T1593.001",
            "T1592.002",
            "T1593",
            "T1669",
            "T1203",
            "T1546",
            "T1127.002",
            "T1647",
            "T1606.001",
            "T1557.003",
            "T1074",
            "T1074.001",
            "T1102.001",
            "T1567.003"
        ]
    }
}