# CTI_reports/report_692_20251222_173448_ground_truth.json
{
    "report_id": "report_692",
    "generated_at": "20251222_173448",
    "scenario_target_ids": [
        "T1584.006",
        "T1583.002",
        "T1659",
        "T1078.004",
        "T1059.001",
        "T1556",
        "T1547.003",
        "T1546.010",
        "T1056.001",
        "T1652",
        "T1021.002",
        "T1090.002"
    ],
    "scenario_flow": [
        "resource-development",
        "resource-development",
        "initial-access",
        "initial-access",
        "execution",
        "persistence",
        "persistence",
        "privilege-escalation",
        "credential-access",
        "discovery",
        "lateral-movement",
        "command-and-control"
    ],
    "llm_declared_ids": [
        "T1584.006",
        "T1583.002",
        "T1659",
        "T1078.004",
        "T1059.001",
        "T1556",
        "T1547.003",
        "T1546.010",
        "T1056.001",
        "T1652",
        "T1021.002",
        "T1090.002"
    ],
    "actual_found_ids": [
        "T1652",
        "T1056.001",
        "T1589",
        "T1659",
        "T1546.010",
        "T1584.006",
        "T1583.002",
        "T1078.004",
        "T1556",
        "T1547.003",
        "T1021.002",
        "T1090.002",
        "T1059.001"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight tactics, techniques, and procedures (TTPs) associated with the fictional threat actor 'ShadowWeaver.' This campaign, observed between January and March 2024, demonstrates a sophisticated multi-stage attack leveraging both cloud and on-premises environments. Initial analysis suggests the actors likely conduct reconnaissance through open-source intelligence gathering, though specific indicators are not detailed in this report.\n\nIn the initial access phase, the actors compromised a GitHub organization to inject malicious code into popular open-source libraries, which was then distributed to victim systems via automated updates. Concurrently, they deployed adversary-in-the-middle DNS servers to impersonate legitimate services, such as corporate login portals, and used content injection through a compromised ISP channel to redirect user traffic through a malicious proxy, facilitating credential capture. For example, logs from a victim network show DNS queries to 'login.corp.example.com' resolving to 192.0.2.1, an unauthorized IP, with subsequent HTTP requests being intercepted.\n\nUpon execution, the adversary employed PowerShell to actively download malicious payloads from command-and-control (C2) servers while simultaneously modifying registry keys to establish persistence. They leveraged compromised cloud accounts, using stolen credentials from victim users to authenticate to Azure tenants and access sensitive resources, such as storage blobs containing financial data. Internal telemetry indicates multiple failed login attempts from IP 203.0.113.5, followed by successful access using valid credentials.\n\nFor persistence and privilege escalation, the actors modified the authentication process by copying the SAM database file from a Volume Shadow Copy and cracking its hashes offline to extract credentials. They also altered time providers via a registry modification to HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\W32Time\\TimeProviders\\, adding a new malicious entry, and configured AppInit DLLs using 'reg add' to ensure execution upon system login. Additionally, they performed keylogging to capture user inputs, including passwords and sensitive data, and conducted device driver discovery using the Windows API function EnumDeviceDrivers() to programmatically retrieve a list of all loaded device drivers for evasion purposes.\n\nIn the lateral movement phase, the actors spread via SMB/Windows Admin Shares, encrypting files on different hosts and potentially compromising the entire network. They employed external proxies to encrypt and redirect C2 traffic, making detection more difficult. For instance, network traffic analysis revealed encrypted connections to 198.51.100.10 via port 443, with payloads masquerading as legitimate HTTPS traffic. The company background includes a mid-sized fintech firm with hybrid infrastructure, but this detail is not directly related to the attack techniques. Recommendations include monitoring for unusual PowerShell scripts, enforcing strong authentication controls, and reviewing DNS configurations for anomalies.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight tactics, techniques, and procedures (TTPs) associated with the fictional threat actor 'ShadowWeaver.' This campaign, observed between January and March 2024, demonstrates a sophisticated multi-stage attack leveraging both cloud and on-premises environments. Initial analysis suggests the actors likely conduct reconnaissance [T1589] through open-source intelligence gathering, though specific indicators are not detailed in this report.\n\nIn the initial access phase, the actors compromised a GitHub organization [T1584.006] to inject malicious code into popular open-source libraries, which was then distributed to victim systems via automated updates. Concurrently, they deployed adversary-in-the-middle DNS servers [T1583.002] to impersonate legitimate services, such as corporate login portals, and used content injection [T1659] through a compromised ISP channel to redirect user traffic through a malicious proxy, facilitating credential capture. For example, logs from a victim network show DNS queries to 'login.corp.example.com' resolving to 192.0.2.1, an unauthorized IP, with subsequent HTTP requests being intercepted.\n\nUpon execution, the adversary employed PowerShell [T1059.001] to actively download malicious payloads from command-and-control (C2) servers while simultaneously modifying registry keys to establish persistence. They leveraged compromised cloud accounts [T1078.004], using stolen credentials from victim users to authenticate to Azure tenants and access sensitive resources, such as storage blobs containing financial data. Internal telemetry indicates multiple failed login attempts from IP 203.0.113.5, followed by successful access using valid credentials.\n\nFor persistence and privilege escalation, the actors modified the authentication process [T1556] by copying the SAM database file from a Volume Shadow Copy and cracking its hashes offline to extract credentials. They also altered time providers [T1547.003] via a registry modification to HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\W32Time\\TimeProviders\\, adding a new malicious entry, and configured AppInit DLLs [T1546.010] using 'reg add' to ensure execution upon system login. Additionally, they performed keylogging [T1056.001] to capture user inputs, including passwords and sensitive data, and conducted device driver discovery [T1652] using the Windows API function EnumDeviceDrivers() to programmatically retrieve a list of all loaded device drivers for evasion purposes.\n\nIn the lateral movement phase, the actors spread via SMB/Windows Admin Shares [T1021.002], encrypting files on different hosts and potentially compromising the entire network. They employed external proxies [T1090.002] to encrypt and redirect C2 traffic, making detection more difficult. For instance, network traffic analysis revealed encrypted connections to 198.51.100.10 via port 443, with payloads masquerading as legitimate HTTPS traffic. The company background includes a mid-sized fintech firm with hybrid infrastructure, but this detail is not directly related to the attack techniques. Recommendations include monitoring for unusual PowerShell scripts, enforcing strong authentication controls, and reviewing DNS configurations for anomalies.",
    "sentence_annotations": [
        {
            "sentence": "Initial analysis suggests the actors likely conduct reconnaissance through open-source intelligence gathering, though specific indicators are not detailed in this report.",
            "technique_id": "T1589",
            "technique_name": "Gather Victim Identity Information"
        },
        {
            "sentence": "002] to impersonate legitimate services, such as corporate login portals, and used content injection through a compromised ISP channel to redirect user traffic through a malicious proxy, facilitating credential capture.",
            "technique_id": "T1659",
            "technique_name": "Content Injection"
        },
        {
            "sentence": "For persistence and privilege escalation, the actors modified the authentication process by copying the SAM database file from a Volume Shadow Copy and cracking its hashes offline to extract credentials.",
            "technique_id": "T1556",
            "technique_name": "Modify Authentication Process"
        },
        {
            "sentence": "001] to capture user inputs, including passwords and sensitive data, and conducted device driver discovery using the Windows API function EnumDeviceDrivers() to programmatically retrieve a list of all loaded device drivers for evasion purposes.",
            "technique_id": "T1652",
            "technique_name": "Device Driver Discovery"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional Campaign 'ShadowWeaver' Targeting Cloud and On-Premises Infrastructure",
        "summary": "This advisory details a simulated campaign by a fictional threat actor, designated as 'ShadowWeaver,' targeting organizations in the technology and finance sectors. The campaign involves supply chain compromise, credential theft, and lateral movement to achieve persistence and data exfiltration. Blue teams are advised to review logs, implement multi-factor authentication, and monitor for suspicious PowerShell activity.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight tactics, techniques, and procedures (TTPs) associated with the fictional threat actor 'ShadowWeaver.' This campaign, observed between January and March 2024, demonstrates a sophisticated multi-stage attack leveraging both cloud and on-premises environments. Initial analysis suggests the actors likely conduct reconnaissance [T1589] through open-source intelligence gathering, though specific indicators are not detailed in this report.\n\nIn the initial access phase, the actors compromised a GitHub organization [T1584.006] to inject malicious code into popular open-source libraries, which was then distributed to victim systems via automated updates. Concurrently, they deployed adversary-in-the-middle DNS servers [T1583.002] to impersonate legitimate services, such as corporate login portals, and used content injection [T1659] through a compromised ISP channel to redirect user traffic through a malicious proxy, facilitating credential capture. For example, logs from a victim network show DNS queries to 'login.corp.example.com' resolving to 192.0.2.1, an unauthorized IP, with subsequent HTTP requests being intercepted.\n\nUpon execution, the adversary employed PowerShell [T1059.001] to actively download malicious payloads from command-and-control (C2) servers while simultaneously modifying registry keys to establish persistence. They leveraged compromised cloud accounts [T1078.004], using stolen credentials from victim users to authenticate to Azure tenants and access sensitive resources, such as storage blobs containing financial data. Internal telemetry indicates multiple failed login attempts from IP 203.0.113.5, followed by successful access using valid credentials.\n\nFor persistence and privilege escalation, the actors modified the authentication process [T1556] by copying the SAM database file from a Volume Shadow Copy and cracking its hashes offline to extract credentials. They also altered time providers [T1547.003] via a registry modification to HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\W32Time\\TimeProviders\\, adding a new malicious entry, and configured AppInit DLLs [T1546.010] using 'reg add' to ensure execution upon system login. Additionally, they performed keylogging [T1056.001] to capture user inputs, including passwords and sensitive data, and conducted device driver discovery [T1652] using the Windows API function EnumDeviceDrivers() to programmatically retrieve a list of all loaded device drivers for evasion purposes.\n\nIn the lateral movement phase, the actors spread via SMB/Windows Admin Shares [T1021.002], encrypting files on different hosts and potentially compromising the entire network. They employed external proxies [T1090.002] to encrypt and redirect C2 traffic, making detection more difficult. For instance, network traffic analysis revealed encrypted connections to 198.51.100.10 via port 443, with payloads masquerading as legitimate HTTPS traffic. The company background includes a mid-sized fintech firm with hybrid infrastructure, but this detail is not directly related to the attack techniques. Recommendations include monitoring for unusual PowerShell scripts, enforcing strong authentication controls, and reviewing DNS configurations for anomalies.",
        "techniques_used": [
            "T1584.006",
            "T1583.002",
            "T1659",
            "T1078.004",
            "T1059.001",
            "T1556",
            "T1547.003",
            "T1546.010",
            "T1056.001",
            "T1652",
            "T1021.002",
            "T1090.002"
        ]
    }
}