# CTI_reports/report_74_20251222_131627_ground_truth.json
{
    "report_id": "report_74",
    "generated_at": "20251222_131627",
    "scenario_target_ids": [
        "T1586.003",
        "T1587.001",
        "T1566.004",
        "T1053.003",
        "T1550.001",
        "T1568.003"
    ],
    "scenario_flow": [
        "resource-development",
        "resource-development",
        "initial-access",
        "execution",
        "lateral-movement",
        "command-and-control"
    ],
    "llm_declared_ids": [
        "T1566.004",
        "T1586.003",
        "T1587.001",
        "T1053.003",
        "T1550.001",
        "T1568.003"
    ],
    "actual_found_ids": [
        "T1566.004",
        "T1587.001",
        "T1586.003",
        "T1550.001",
        "T1568.003",
        "T1053.003"
    ],
    "clean_text": "Background: CLEAVER, a fictional advanced persistent threat (APT) group, has been observed targeting financial and technology sectors in a simulated exercise. The group employs sophisticated tactics to infiltrate cloud-based infrastructures, leveraging social engineering and custom tools. Initial analysis indicates the campaign began in early 2024, with no real-world impact reported. Logs from affected systems show increased network traffic and unauthorized access attempts, but these are part of the training scenario.\n\nAttack Narrative: The actors initiated the campaign by conducting spear-phishing via voice calls, directing employees to download a remote access trojan disguised as a document viewer, while simultaneously sending phishing emails disguised as IT support to harvest cloud account credentials. Upon execution, CLEAVER deployed custom malware that performed ARP poisoning and keystroke logging, and used crontabs to establish persistence. The malware then exploited OAuth tokens to bypass multi-factor authentication and access cloud resources, while tunneling TCP traffic through DNS queries using iodine to evade network filtering. System logs from a fictional company, CloudSecure Inc., recorded anomalous DNS requests and cron job modifications, though no data exfiltration occurred in this training scenario. The actors maintained access for several days, simulating lateral movement within the environment before being detected by security teams.",
    "tagged_text": "Background: CLEAVER, a fictional advanced persistent threat (APT) group, has been observed targeting financial and technology sectors in a simulated exercise. The group employs sophisticated tactics to infiltrate cloud-based infrastructures, leveraging social engineering and custom tools. Initial analysis indicates the campaign began in early 2024, with no real-world impact reported. Logs from affected systems show increased network traffic and unauthorized access attempts, but these are part of the training scenario.\n\nAttack Narrative: The actors initiated the campaign by conducting spear-phishing via voice calls [T1566.004], directing employees to download a remote access trojan disguised as a document viewer, while simultaneously sending phishing emails disguised as IT support to harvest cloud account credentials [T1586.003]. Upon execution, CLEAVER deployed custom malware [T1587.001] that performed ARP poisoning and keystroke logging, and used crontabs to establish persistence [T1053.003]. The malware then exploited OAuth tokens to bypass multi-factor authentication and access cloud resources [T1550.001], while tunneling TCP traffic through DNS queries using iodine to evade network filtering [T1568.003]. System logs from a fictional company, CloudSecure Inc., recorded anomalous DNS requests and cron job modifications, though no data exfiltration occurred in this training scenario. The actors maintained access for several days, simulating lateral movement within the environment before being detected by security teams.",
    "sentence_annotations": [],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional Campaign by Threat Actor CLEAVER",
        "summary": "This advisory details a simulated campaign by the fictional threat actor CLEAVER, targeting cloud environments through multi-vector phishing, custom malware deployment, and DNS tunneling. The campaign involves credential theft, persistence via cron jobs, and token abuse to bypass security controls, designed for blue team training to enhance detection and response capabilities.",
        "full_report_text": "Background: CLEAVER, a fictional advanced persistent threat (APT) group, has been observed targeting financial and technology sectors in a simulated exercise. The group employs sophisticated tactics to infiltrate cloud-based infrastructures, leveraging social engineering and custom tools. Initial analysis indicates the campaign began in early 2024, with no real-world impact reported. Logs from affected systems show increased network traffic and unauthorized access attempts, but these are part of the training scenario.\n\nAttack Narrative: The actors initiated the campaign by conducting spear-phishing via voice calls [T1566.004], directing employees to download a remote access trojan disguised as a document viewer, while simultaneously sending phishing emails disguised as IT support to harvest cloud account credentials [T1586.003]. Upon execution, CLEAVER deployed custom malware [T1587.001] that performed ARP poisoning and keystroke logging, and used crontabs to establish persistence [T1053.003]. The malware then exploited OAuth tokens to bypass multi-factor authentication and access cloud resources [T1550.001], while tunneling TCP traffic through DNS queries using iodine to evade network filtering [T1568.003]. System logs from a fictional company, CloudSecure Inc., recorded anomalous DNS requests and cron job modifications, though no data exfiltration occurred in this training scenario. The actors maintained access for several days, simulating lateral movement within the environment before being detected by security teams.",
        "techniques_used": [
            "T1566.004",
            "T1586.003",
            "T1587.001",
            "T1053.003",
            "T1550.001",
            "T1568.003"
        ]
    }
}