# CTI_reports/report_90_20251222_115259_ground_truth.json
{
    "report_id": "report_90",
    "generated_at": "20251222_115259",
    "scenario_target_ids": [
        "T1583.006",
        "T1566.001",
        "T1200",
        "T1053",
        "T1547.008",
        "T1098.003",
        "T1036.005",
        "T1563",
        "T1080",
        "T1090.004",
        "T1568",
        "T1029"
    ],
    "scenario_flow": [
        "resource-development",
        "initial-access",
        "initial-access",
        "execution",
        "persistence",
        "privilege-escalation",
        "defense-evasion",
        "lateral-movement",
        "lateral-movement",
        "command-and-control",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1589",
        "T1583.006",
        "T1566.001",
        "T1053",
        "T1547.008",
        "T1098.003",
        "T1036.005",
        "T1090.004",
        "T1080",
        "T1563",
        "T1568",
        "T1029",
        "T1200"
    ],
    "actual_found_ids": [
        "T1053",
        "T1036.005",
        "T1029",
        "T1566.001",
        "T1547.008",
        "T1568",
        "T1589",
        "T1090.004",
        "T1098.003",
        "T1200",
        "T1080",
        "T1563",
        "T1583.006"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to inform blue teams about a fictional APT campaign, 'ShadowWeave', observed targeting energy and financial sectors. This report is based on a simulated attack scenario for educational purposes. Initial Access: The actors likely conducted reconnaissance and abused web services by sending spearphishing messages via social media platforms, directing victims to credential harvesting sites. Execution: Upon user interaction, the campaign delivered spearphishing attachments, including malicious PDFs with embedded links and archive files containing executable payloads, which were executed to gain initial foothold. Persistence: To maintain access, the adversary employed scheduled tasks using schtasks.exe to run malicious droppers at user logon, while simultaneously replacing a legitimate LSASS driver with a backdoored version and modifying cloud permissions to grant administrative privileges for stealthy persistence. Defense Evasion: The actors used match legitimate resource names by renaming tools to appear as benign programs and employed domain fronting by altering TLS header fields to bypass network filters, while also tainting shared content by writing malicious executables to network shares accessible by multiple users. Lateral Movement: The threat actor hijacked remote service sessions using stolen session tokens to access RDP connections and exfiltrate sensitive data. Command and Control: For covert communication, the malware performed dynamic resolution through DNS queries to establish command and control channels. Exfiltration: Data was exfiltrated via scheduled transfers by writing PowerShell scripts to temporary directories for automated data extraction. Additionally, system logs from compromised servers showed unusual network traffic spikes, and company background checks revealed no prior incidents, though internal audits flagged unauthorized hardware additions, such as a passive tap discovered on an Ethernet line in the server room, which was unrelated to the primary campaign but highlighted physical security gaps.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to inform blue teams about a fictional APT campaign, 'ShadowWeave', observed targeting energy and financial sectors. This report is based on a simulated attack scenario for educational purposes. Initial Access: The actors likely conducted reconnaissance [T1589] and abused web services [T1583.006] by sending spearphishing messages via social media platforms, directing victims to credential harvesting sites. Execution: Upon user interaction, the campaign delivered spearphishing attachments [T1566.001], including malicious PDFs with embedded links and archive files containing executable payloads, which were executed to gain initial foothold. Persistence: To maintain access, the adversary employed scheduled tasks [T1053] using schtasks.exe to run malicious droppers at user logon, while simultaneously replacing a legitimate LSASS driver with a backdoored version [T1547.008] and modifying cloud permissions [T1098.003] to grant administrative privileges for stealthy persistence. Defense Evasion: The actors used match legitimate resource names [T1036.005] by renaming tools to appear as benign programs and employed domain fronting [T1090.004] by altering TLS header fields to bypass network filters, while also tainting shared content [T1080] by writing malicious executables to network shares accessible by multiple users. Lateral Movement: The threat actor hijacked remote service sessions [T1563] using stolen session tokens to access RDP connections and exfiltrate sensitive data. Command and Control: For covert communication, the malware performed dynamic resolution [T1568] through DNS queries to establish command and control channels. Exfiltration: Data was exfiltrated via scheduled transfers [T1029] by writing PowerShell scripts to temporary directories for automated data extraction. Additionally, system logs from compromised servers showed unusual network traffic spikes, and company background checks revealed no prior incidents, though internal audits flagged unauthorized hardware additions [T1200], such as a passive tap discovered on an Ethernet line in the server room, which was unrelated to the primary campaign but highlighted physical security gaps.",
    "sentence_annotations": [
        {
            "sentence": "Initial Access: The actors likely conducted reconnaissance and abused web services [T1583.",
            "technique_id": "T1589",
            "technique_name": "Gather Victim Identity Information"
        },
        {
            "sentence": "Persistence: To maintain access, the adversary employed scheduled tasks using schtasks.",
            "technique_id": "T1053",
            "technique_name": "Scheduled Task/Job"
        },
        {
            "sentence": "004] by altering TLS header fields to bypass network filters, while also tainting shared content by writing malicious executables to network shares accessible by multiple users.",
            "technique_id": "T1080",
            "technique_name": "Taint Shared Content"
        },
        {
            "sentence": "Lateral Movement: The threat actor hijacked remote service sessions using stolen session tokens to access RDP connections and exfiltrate sensitive data.",
            "technique_id": "T1563",
            "technique_name": "Remote Service Session Hijacking"
        },
        {
            "sentence": "Command and Control: For covert communication, the malware performed dynamic resolution through DNS queries to establish command and control channels.",
            "technique_id": "T1568",
            "technique_name": "Dynamic Resolution"
        },
        {
            "sentence": "Exfiltration: Data was exfiltrated via scheduled transfers by writing PowerShell scripts to temporary directories for automated data extraction.",
            "technique_id": "T1029",
            "technique_name": "Scheduled Transfer"
        },
        {
            "sentence": "Additionally, system logs from compromised servers showed unusual network traffic spikes, and company background checks revealed no prior incidents, though internal audits flagged unauthorized hardware additions, such as a passive tap discovered on an Ethernet line in the server room, which was unrelated to the primary campaign but highlighted physical security gaps.",
            "technique_id": "T1200",
            "technique_name": "Hardware Additions"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'ShadowWeave' Targeting Critical Infrastructure",
        "summary": "This advisory details a simulated campaign by a fictional advanced persistent threat (APT) group, designated 'ShadowWeave', targeting critical infrastructure sectors. The campaign employs sophisticated techniques across the kill chain, from initial access via social engineering to persistence and exfiltration. Blue teams should use this report for training to enhance detection and response capabilities.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to inform blue teams about a fictional APT campaign, 'ShadowWeave', observed targeting energy and financial sectors. This report is based on a simulated attack scenario for educational purposes. Initial Access: The actors likely conducted reconnaissance [T1589] and abused web services [T1583.006] by sending spearphishing messages via social media platforms, directing victims to credential harvesting sites. Execution: Upon user interaction, the campaign delivered spearphishing attachments [T1566.001], including malicious PDFs with embedded links and archive files containing executable payloads, which were executed to gain initial foothold. Persistence: To maintain access, the adversary employed scheduled tasks [T1053] using schtasks.exe to run malicious droppers at user logon, while simultaneously replacing a legitimate LSASS driver with a backdoored version [T1547.008] and modifying cloud permissions [T1098.003] to grant administrative privileges for stealthy persistence. Defense Evasion: The actors used match legitimate resource names [T1036.005] by renaming tools to appear as benign programs and employed domain fronting [T1090.004] by altering TLS header fields to bypass network filters, while also tainting shared content [T1080] by writing malicious executables to network shares accessible by multiple users. Lateral Movement: The threat actor hijacked remote service sessions [T1563] using stolen session tokens to access RDP connections and exfiltrate sensitive data. Command and Control: For covert communication, the malware performed dynamic resolution [T1568] through DNS queries to establish command and control channels. Exfiltration: Data was exfiltrated via scheduled transfers [T1029] by writing PowerShell scripts to temporary directories for automated data extraction. Additionally, system logs from compromised servers showed unusual network traffic spikes, and company background checks revealed no prior incidents, though internal audits flagged unauthorized hardware additions [T1200], such as a passive tap discovered on an Ethernet line in the server room, which was unrelated to the primary campaign but highlighted physical security gaps.",
        "techniques_used": [
            "T1589",
            "T1583.006",
            "T1566.001",
            "T1053",
            "T1547.008",
            "T1098.003",
            "T1036.005",
            "T1090.004",
            "T1080",
            "T1563",
            "T1568",
            "T1029",
            "T1200"
        ]
    }
}