# CTI_reports/report_370_20251222_151835_ground_truth.json
{
    "report_id": "report_370",
    "generated_at": "20251222_151835",
    "scenario_target_ids": [
        "T1584",
        "T1650",
        "T1078.002",
        "T1053.007",
        "T1574.011",
        "T1021.006",
        "T1090.004"
    ],
    "scenario_flow": [
        "resource-development",
        "resource-development",
        "initial-access",
        "execution",
        "privilege-escalation",
        "lateral-movement",
        "command-and-control"
    ],
    "llm_declared_ids": [
        "T1584",
        "T1650",
        "T1078.002",
        "T1053.007",
        "T1574.011",
        "T1021.006",
        "T1090.004"
    ],
    "actual_found_ids": [
        "T1021.006",
        "T1078.002",
        "T1053.007",
        "T1584",
        "T1090.004",
        "T1574.011",
        "T1650"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight a fictional APT campaign designed for educational purposes. The threat actors initiated their operation by compromising a third-party web server, which hosted a malicious executable (malicious.exe) that was later used for initial access. Upon gaining entry, the attacker uploaded a PHP web shell to a vulnerable web server to establish and maintain persistent access. System logs from the victim organization, a mid-sized financial services company, showed unusual network traffic spikes during this phase, though initial firewall alerts were dismissed as false positives. The actors then leveraged stolen credentials to authenticate using valid domain accounts, enabling them to move laterally within the network. As part of their lateral movement strategy, they configured a CronJob to execute a container designed for malicious activities every 6 hours, while simultaneously hijacking the Messenger service registry key in HKLM\\SYSTEM\\CurrentControlSet\\Services to ensure persistence. During this process, internal monitoring tools detected anomalous registry modifications, but no immediate action was taken due to resource constraints. The threat actor employed Windows Remote Management (WinRM) to remotely deploy additional payloads across enterprise systems, and used Cobalt Strike with domain fronting capabilities by manipulating HTTP Host Headers to evade network-based detection. Background information on the company's recent infrastructure upgrades, including migration to cloud services, was irrelevant to the attack but noted in internal audits. This campaign underscores the importance of monitoring for multi-vector attacks and implementing robust security controls.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight a fictional APT campaign designed for educational purposes. The threat actors initiated their operation by compromising a third-party web server [T1584], which hosted a malicious executable (malicious.exe) that was later used for initial access. Upon gaining entry, the attacker uploaded a PHP web shell to a vulnerable web server to establish and maintain persistent access [T1650]. System logs from the victim organization, a mid-sized financial services company, showed unusual network traffic spikes during this phase, though initial firewall alerts were dismissed as false positives. The actors then leveraged stolen credentials to authenticate using valid domain accounts [T1078.002], enabling them to move laterally within the network. As part of their lateral movement strategy, they configured a CronJob to execute a container designed for malicious activities every 6 hours [T1053.007], while simultaneously hijacking the Messenger service registry key in HKLM\\SYSTEM\\CurrentControlSet\\Services to ensure persistence [T1574.011]. During this process, internal monitoring tools detected anomalous registry modifications, but no immediate action was taken due to resource constraints. The threat actor employed Windows Remote Management (WinRM) to remotely deploy additional payloads across enterprise systems [T1021.006], and used Cobalt Strike with domain fronting capabilities by manipulating HTTP Host Headers to evade network-based detection [T1090.004]. Background information on the company's recent infrastructure upgrades, including migration to cloud services, was irrelevant to the attack but noted in internal audits. This campaign underscores the importance of monitoring for multi-vector attacks and implementing robust security controls.",
    "sentence_annotations": [
        {
            "sentence": "The threat actors initiated their operation by compromising a third-party web server, which hosted a malicious executable (malicious.",
            "technique_id": "T1584",
            "technique_name": "Compromise Infrastructure"
        },
        {
            "sentence": "Upon gaining entry, the attacker uploaded a PHP web shell to a vulnerable web server to establish and maintain persistent access.",
            "technique_id": "T1650",
            "technique_name": "Acquire Access"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign Targeting Enterprise Infrastructure",
        "summary": "This advisory details a simulated campaign by a fictional advanced persistent threat (APT) group targeting enterprise networks. The campaign involves compromised infrastructure, web shell deployment, credential abuse, container orchestration, persistence mechanisms, remote management, and domain fronting for evasion. This report is intended for blue team training to enhance detection and response capabilities.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight a fictional APT campaign designed for educational purposes. The threat actors initiated their operation by compromising a third-party web server [T1584], which hosted a malicious executable (malicious.exe) that was later used for initial access. Upon gaining entry, the attacker uploaded a PHP web shell to a vulnerable web server to establish and maintain persistent access [T1650]. System logs from the victim organization, a mid-sized financial services company, showed unusual network traffic spikes during this phase, though initial firewall alerts were dismissed as false positives. The actors then leveraged stolen credentials to authenticate using valid domain accounts [T1078.002], enabling them to move laterally within the network. As part of their lateral movement strategy, they configured a CronJob to execute a container designed for malicious activities every 6 hours [T1053.007], while simultaneously hijacking the Messenger service registry key in HKLM\\SYSTEM\\CurrentControlSet\\Services to ensure persistence [T1574.011]. During this process, internal monitoring tools detected anomalous registry modifications, but no immediate action was taken due to resource constraints. The threat actor employed Windows Remote Management (WinRM) to remotely deploy additional payloads across enterprise systems [T1021.006], and used Cobalt Strike with domain fronting capabilities by manipulating HTTP Host Headers to evade network-based detection [T1090.004]. Background information on the company's recent infrastructure upgrades, including migration to cloud services, was irrelevant to the attack but noted in internal audits. This campaign underscores the importance of monitoring for multi-vector attacks and implementing robust security controls.",
        "techniques_used": [
            "T1584",
            "T1650",
            "T1078.002",
            "T1053.007",
            "T1574.011",
            "T1021.006",
            "T1090.004"
        ]
    }
}