# CTI_reports/report_659_20251222_172144_ground_truth.json
{
    "report_id": "report_659",
    "generated_at": "20251222_172144",
    "scenario_target_ids": [
        "T1586.002",
        "T1585",
        "T1583.006",
        "T1195.001",
        "T1059.011",
        "T1037",
        "T1078.004",
        "T1550",
        "T1218.007",
        "T1552",
        "T1021.006",
        "T1602",
        "T1102",
        "T1011.001"
    ],
    "scenario_flow": [
        "resource-development",
        "resource-development",
        "resource-development",
        "initial-access",
        "execution",
        "persistence",
        "persistence",
        "defense-evasion",
        "defense-evasion",
        "credential-access",
        "lateral-movement",
        "collection",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1586.002",
        "T1585",
        "T1583.006",
        "T1195.001",
        "T1059.011",
        "T1037",
        "T1078.004",
        "T1550",
        "T1218.007",
        "T1552",
        "T1021.006",
        "T1602",
        "T1102",
        "T1011.001"
    ],
    "actual_found_ids": [
        "T1011.001",
        "T1218.007",
        "T1059.011",
        "T1585",
        "T1583.006",
        "T1195.001",
        "T1550",
        "T1552",
        "T1021.006",
        "T1602",
        "T1078.004",
        "T1586.002",
        "T1102",
        "T1037"
    ],
    "clean_text": "On April 15, 2024, multiple financial organizations reported suspicious activity linked to the fictional threat actor ShadowStrike, known for targeting financial records and intellectual property. The campaign began with reconnaissance where the actor established fake accounts on social media platforms, including a fabricated Facebook profile using stolen personal data, and abused Twitter's API to create multiple accounts for distributing phishing links. Initial access was achieved through phishing emails sent from compromised email accounts, which contained malicious attachments. Once inside the network, the actor executed malicious code by injecting it into a deployment script deploy.ps1 to capture financial records, while simultaneously using Lua modules for execution and writing a malicious VBScript to the scheduled tasks directory for automatic execution upon system boot. Persistence and lateral movement involved exploiting cloud APIs to create backdoor access to virtual machine instances, performing a Pass-the-Hash attack using stolen NTLM hashes to authenticate to remote systems, and using msiexec.exe to invoke a malicious DLL disguised as a legitimate component, bypassing application controls. The actor also obtained credentials insecurely stored on network appliances and used winrm commands to remotely modify firewall settings on Windows servers. Data collection included extracting configuration data from the local Group Policy repository to understand domain security settings, and the actor leveraged a legitimate web service for evading detection during command and control. Exfiltration was conducted by exploiting the Bluetooth service discovery protocol to transfer stolen intellectual property. System logs from affected organizations showed routine network traffic spikes unrelated to the attack, and background checks revealed no prior incidents involving ShadowStrike, though the actor's tactics align with historical APT patterns in the financial sector.",
    "tagged_text": "On April 15, 2024, multiple financial organizations reported suspicious activity linked to the fictional threat actor ShadowStrike, known for targeting financial records and intellectual property. The campaign began with reconnaissance where the actor established fake accounts on social media platforms [T1585], including a fabricated Facebook profile using stolen personal data, and abused Twitter's API to create multiple accounts for distributing phishing links [T1583.006]. Initial access was achieved through phishing emails sent from compromised email accounts [T1586.002], which contained malicious attachments. Once inside the network, the actor executed malicious code by injecting it into a deployment script deploy.ps1 to capture financial records [T1195.001], while simultaneously using Lua modules for execution [T1059.011] and writing a malicious VBScript to the scheduled tasks directory for automatic execution upon system boot [T1037]. Persistence and lateral movement involved exploiting cloud APIs to create backdoor access to virtual machine instances [T1078.004], performing a Pass-the-Hash attack using stolen NTLM hashes to authenticate to remote systems [T1550], and using msiexec.exe to invoke a malicious DLL disguised as a legitimate component, bypassing application controls [T1218.007]. The actor also obtained credentials insecurely stored on network appliances [T1552] and used winrm commands to remotely modify firewall settings on Windows servers [T1021.006]. Data collection included extracting configuration data from the local Group Policy repository to understand domain security settings [T1602], and the actor leveraged a legitimate web service for evading detection during command and control [T1102]. Exfiltration was conducted by exploiting the Bluetooth service discovery protocol to transfer stolen intellectual property [T1011.001]. System logs from affected organizations showed routine network traffic spikes unrelated to the attack, and background checks revealed no prior incidents involving ShadowStrike, though the actor's tactics align with historical APT patterns in the financial sector.",
    "sentence_annotations": [
        {
            "sentence": "The campaign began with reconnaissance where the actor established fake accounts on social media platforms, including a fabricated Facebook profile using stolen personal data, and abused Twitter's API to create multiple accounts for distributing phishing links [T1583.",
            "technique_id": "T1585",
            "technique_name": "Establish Accounts"
        },
        {
            "sentence": "011] and writing a malicious VBScript to the scheduled tasks directory for automatic execution upon system boot.",
            "technique_id": "T1037",
            "technique_name": "Boot or Logon Initialization Scripts"
        },
        {
            "sentence": "004], performing a Pass-the-Hash attack using stolen NTLM hashes to authenticate to remote systems, and using msiexec.",
            "technique_id": "T1550",
            "technique_name": "Use Alternate Authentication Material"
        },
        {
            "sentence": "The actor also obtained credentials insecurely stored on network appliances and used winrm commands to remotely modify firewall settings on Windows servers [T1021.",
            "technique_id": "T1552",
            "technique_name": "Unsecured Credentials"
        },
        {
            "sentence": "Data collection included extracting configuration data from the local Group Policy repository to understand domain security settings, and the actor leveraged a legitimate web service for evading detection during command and control.",
            "technique_id": "T1602",
            "technique_name": "Data from Configuration Repository"
        },
        {
            "sentence": "Data collection included extracting configuration data from the local Group Policy repository to understand domain security settings, and the actor leveraged a legitimate web service for evading detection during command and control.",
            "technique_id": "T1102",
            "technique_name": "Web Service"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'ShadowStrike' Targeting Financial Sector",
        "summary": "This advisory details a simulated campaign by the fictional threat actor 'ShadowStrike' targeting financial institutions, using a blend of social engineering, supply chain compromise, and cloud exploitation to exfiltrate sensitive data. The activity involves techniques such as establishing fake accounts, abusing web services, and leveraging alternate authentication methods, with exfiltration over Bluetooth observed in later stages.",
        "full_report_text": "On April 15, 2024, multiple financial organizations reported suspicious activity linked to the fictional threat actor ShadowStrike, known for targeting financial records and intellectual property. The campaign began with reconnaissance where the actor established fake accounts on social media platforms [T1585], including a fabricated Facebook profile using stolen personal data, and abused Twitter's API to create multiple accounts for distributing phishing links [T1583.006]. Initial access was achieved through phishing emails sent from compromised email accounts [T1586.002], which contained malicious attachments. Once inside the network, the actor executed malicious code by injecting it into a deployment script deploy.ps1 to capture financial records [T1195.001], while simultaneously using Lua modules for execution [T1059.011] and writing a malicious VBScript to the scheduled tasks directory for automatic execution upon system boot [T1037]. Persistence and lateral movement involved exploiting cloud APIs to create backdoor access to virtual machine instances [T1078.004], performing a Pass-the-Hash attack using stolen NTLM hashes to authenticate to remote systems [T1550], and using msiexec.exe to invoke a malicious DLL disguised as a legitimate component, bypassing application controls [T1218.007]. The actor also obtained credentials insecurely stored on network appliances [T1552] and used winrm commands to remotely modify firewall settings on Windows servers [T1021.006]. Data collection included extracting configuration data from the local Group Policy repository to understand domain security settings [T1602], and the actor leveraged a legitimate web service for evading detection during command and control [T1102]. Exfiltration was conducted by exploiting the Bluetooth service discovery protocol to transfer stolen intellectual property [T1011.001]. System logs from affected organizations showed routine network traffic spikes unrelated to the attack, and background checks revealed no prior incidents involving ShadowStrike, though the actor's tactics align with historical APT patterns in the financial sector.",
        "techniques_used": [
            "T1586.002",
            "T1585",
            "T1583.006",
            "T1195.001",
            "T1059.011",
            "T1037",
            "T1078.004",
            "T1550",
            "T1218.007",
            "T1552",
            "T1021.006",
            "T1602",
            "T1102",
            "T1011.001"
        ]
    }
}