# CTI_reports/report_564_20251222_164301_ground_truth.json
{
    "report_id": "report_564",
    "generated_at": "20251222_164301",
    "scenario_target_ids": [
        "T1585.001",
        "T1078.004",
        "T1569.003",
        "T1547",
        "T1547.007",
        "T1003.001",
        "T1557.001",
        "T1021.006",
        "T1568.002",
        "T1011.001"
    ],
    "scenario_flow": [
        "resource-development",
        "initial-access",
        "execution",
        "persistence",
        "privilege-escalation",
        "credential-access",
        "credential-access",
        "lateral-movement",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1585.001",
        "T1078.004",
        "T1569.003",
        "T1547",
        "T1547.007",
        "T1003.001",
        "T1557.001",
        "T1021.006",
        "T1568.002",
        "T1011.001"
    ],
    "actual_found_ids": [
        "T1011.001",
        "T1547.007",
        "T1568.002",
        "T1547",
        "T1585.001",
        "T1557.001",
        "T1021.006",
        "T1078.004",
        "T1569.003",
        "T1003.001"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight the activities of the advanced persistent threat (APT) group 'ShadowLink'. This fictional campaign is designed for educational purposes to train blue teams in detecting and responding to sophisticated cyber threats. Initial reconnaissance involved the creation of fake LinkedIn profiles with detailed connections and photos to gather intelligence on target employees. Subsequently, the actors exploited compromised Azure Managed Identities to gain unauthorized access to virtual machines in the cloud environment, establishing a foothold for further operations. On Linux systems, the adversary executed 'systemctl start network-sniffer' to activate packet capture capabilities for monitoring network traffic, while on Windows systems, they modified the HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run registry key to ensure malicious.exe persisted across reboots. Additionally, on macOS, the threat actor altered com.apple.loginwindow.UUID.plist to automatically launch network scanning utilities upon user login, maintaining stealthy access. For credential harvesting, the group deployed a modified version of Mimikatz to extract credentials from LSASS memory, facilitating lateral movement within the network. They also employed LLMNR/NBT-NS poisoning techniques to intercept plaintext credentials and relay SMB traffic, enabling unauthorized access to other systems. Using Windows Remote Management (WinRM) commands, the actors remotely adjusted firewall settings on compromised Windows machines to open ports for command and control (C2) communication. To reestablish persistent C2 channels, the compromised systems executed 'c2_restore.exe', which utilized domain generation algorithms (DGAs) to dynamically connect to new C2 servers. Finally, data exfiltration was conducted over Bluetooth by exploiting SDP vulnerabilities to establish unauthorized connections and transfer stolen information. Throughout the campaign, system logs showed routine updates and user activity, with no immediate alerts from standard monitoring tools, highlighting the need for enhanced detection mechanisms. The targeted organization, a mid-sized tech firm with hybrid cloud deployments, experienced minor disruptions but no actual data loss in this training scenario.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight the activities of the advanced persistent threat (APT) group 'ShadowLink'. This fictional campaign is designed for educational purposes to train blue teams in detecting and responding to sophisticated cyber threats. Initial reconnaissance involved the creation of fake LinkedIn profiles with detailed connections and photos to gather intelligence on target employees [T1585.001]. Subsequently, the actors exploited compromised Azure Managed Identities to gain unauthorized access to virtual machines in the cloud environment, establishing a foothold for further operations [T1078.004]. On Linux systems, the adversary executed 'systemctl start network-sniffer' to activate packet capture capabilities for monitoring network traffic [T1569.003], while on Windows systems, they modified the HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run registry key to ensure malicious.exe persisted across reboots [T1547]. Additionally, on macOS, the threat actor altered com.apple.loginwindow.UUID.plist to automatically launch network scanning utilities upon user login, maintaining stealthy access [T1547.007]. For credential harvesting, the group deployed a modified version of Mimikatz to extract credentials from LSASS memory, facilitating lateral movement within the network [T1003.001]. They also employed LLMNR/NBT-NS poisoning techniques to intercept plaintext credentials and relay SMB traffic, enabling unauthorized access to other systems [T1557.001]. Using Windows Remote Management (WinRM) commands, the actors remotely adjusted firewall settings on compromised Windows machines to open ports for command and control (C2) communication [T1021.006]. To reestablish persistent C2 channels, the compromised systems executed 'c2_restore.exe', which utilized domain generation algorithms (DGAs) to dynamically connect to new C2 servers [T1568.002]. Finally, data exfiltration was conducted over Bluetooth by exploiting SDP vulnerabilities to establish unauthorized connections and transfer stolen information [T1011.001]. Throughout the campaign, system logs showed routine updates and user activity, with no immediate alerts from standard monitoring tools, highlighting the need for enhanced detection mechanisms. The targeted organization, a mid-sized tech firm with hybrid cloud deployments, experienced minor disruptions but no actual data loss in this training scenario.",
    "sentence_annotations": [
        {
            "sentence": "exe persisted across reboots.",
            "technique_id": "T1547",
            "technique_name": "Boot or Logon Autostart Execution"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'ShadowLink' Targeting Hybrid Cloud Environments",
        "summary": "This advisory details a simulated campaign by the fictional threat actor group 'ShadowLink', which targets organizations with hybrid cloud infrastructures. The campaign involves social engineering, cloud account compromise, and persistent access mechanisms to exfiltrate sensitive data. Blue teams should use this report for training to identify and mitigate similar tactics, techniques, and procedures (TTPs).",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight the activities of the advanced persistent threat (APT) group 'ShadowLink'. This fictional campaign is designed for educational purposes to train blue teams in detecting and responding to sophisticated cyber threats. Initial reconnaissance involved the creation of fake LinkedIn profiles with detailed connections and photos to gather intelligence on target employees [T1585.001]. Subsequently, the actors exploited compromised Azure Managed Identities to gain unauthorized access to virtual machines in the cloud environment, establishing a foothold for further operations [T1078.004]. On Linux systems, the adversary executed 'systemctl start network-sniffer' to activate packet capture capabilities for monitoring network traffic [T1569.003], while on Windows systems, they modified the HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run registry key to ensure malicious.exe persisted across reboots [T1547]. Additionally, on macOS, the threat actor altered com.apple.loginwindow.UUID.plist to automatically launch network scanning utilities upon user login, maintaining stealthy access [T1547.007]. For credential harvesting, the group deployed a modified version of Mimikatz to extract credentials from LSASS memory, facilitating lateral movement within the network [T1003.001]. They also employed LLMNR/NBT-NS poisoning techniques to intercept plaintext credentials and relay SMB traffic, enabling unauthorized access to other systems [T1557.001]. Using Windows Remote Management (WinRM) commands, the actors remotely adjusted firewall settings on compromised Windows machines to open ports for command and control (C2) communication [T1021.006]. To reestablish persistent C2 channels, the compromised systems executed 'c2_restore.exe', which utilized domain generation algorithms (DGAs) to dynamically connect to new C2 servers [T1568.002]. Finally, data exfiltration was conducted over Bluetooth by exploiting SDP vulnerabilities to establish unauthorized connections and transfer stolen information [T1011.001]. Throughout the campaign, system logs showed routine updates and user activity, with no immediate alerts from standard monitoring tools, highlighting the need for enhanced detection mechanisms. The targeted organization, a mid-sized tech firm with hybrid cloud deployments, experienced minor disruptions but no actual data loss in this training scenario.",
        "techniques_used": [
            "T1585.001",
            "T1078.004",
            "T1569.003",
            "T1547",
            "T1547.007",
            "T1003.001",
            "T1557.001",
            "T1021.006",
            "T1568.002",
            "T1011.001"
        ]
    }
}