# CTI_reports/report_779_20251223_111044_ground_truth.json
{
    "report_id": "report_779",
    "generated_at": "20251223_111044",
    "scenario_target_ids": [
        "T1598.003",
        "T1584.004",
        "T1650",
        "T1091",
        "T1053.005",
        "T1546",
        "T1574.006",
        "T1027.007",
        "T1070.009",
        "T1036.007",
        "T1132.001",
        "T1048.003",
        "T1489"
    ],
    "scenario_flow": [
        "reconnaissance",
        "resource-development",
        "resource-development",
        "initial-access",
        "execution",
        "persistence",
        "persistence",
        "defense-evasion",
        "defense-evasion",
        "defense-evasion",
        "command-and-control",
        "exfiltration",
        "impact"
    ],
    "llm_declared_ids": [
        "T1598.003",
        "T1584.004",
        "T1650",
        "T1091",
        "T1053.005",
        "T1546",
        "T1574.006",
        "T1027.007",
        "T1070.009",
        "T1036.007",
        "T1132.001",
        "T1048.003",
        "T1489"
    ],
    "actual_found_ids": [
        "T1053.005",
        "T1546",
        "T1584.004",
        "T1589",
        "T1091",
        "T1132.001",
        "T1574.006",
        "T1489",
        "T1027.007",
        "T1048.003",
        "T1105",
        "T1598.003",
        "T1070.009",
        "T1036.007",
        "T1650",
        "T1059.001"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide blue teams with training on a fictional advanced persistent threat (APT) campaign, designated 'Nightshade'. This scenario is based on observed tactics, techniques, and procedures (TTPs) in real-world incidents but is entirely fabricated for educational purposes. Nightshade actors likely conduct reconnaissance on financial targets, though this phase is not detailed in the scenario. The campaign begins with Initial Access, where the adversary sent spearphishing emails disguised as password reset requests, containing malicious links to compromised web servers used for command and control (C2). To enhance access, they acquired pre-compromised VPN credentials from an access broker network, enabling lateral movement. During Execution, the malware employed PowerShell to download payloads and replicate through removable media for propagation across air-gapped systems. For Persistence, the actors created scheduled tasks and implemented event-triggered execution through COM hijacking, while also modifying /etc/ld.so.preload to hijack dynamic linker functions and obscure behavior via dynamic API resolution using GetProcAddress. In Defense Evasion, they used double file extensions to masquerade malicious LNK files as PDF objects and cleared persistence traces by removing IFEO registry values. For Command and Control, C2 traffic was encoded with Base64 to evade detection. During Exfiltration, sensitive documents were collected and sent via unencrypted HTTP requests to remote servers. Finally, for Impact, the actors issued commands to stop critical services on compromised hosts, disrupting operations. Logs from affected systems show routine network traffic and user login attempts, but no direct correlation to the attack. Company background: The targeted entity is a mid-sized bank with outdated security protocols, though this information is irrelevant to the TTPs. This report serves as a training tool to enhance defensive capabilities against similar real-world threats.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide blue teams with training on a fictional advanced persistent threat (APT) campaign, designated 'Nightshade'. This scenario is based on observed tactics, techniques, and procedures (TTPs) in real-world incidents but is entirely fabricated for educational purposes. Nightshade actors likely conduct reconnaissance [T1589] on financial targets, though this phase is not detailed in the scenario. The campaign begins with Initial Access, where the adversary sent spearphishing emails [T1598.003] disguised as password reset requests, containing malicious links to compromised web servers [T1584.004] used for command and control (C2). To enhance access, they acquired pre-compromised VPN credentials [T1650] from an access broker network, enabling lateral movement. During Execution, the malware employed PowerShell [T1059.001] to download payloads [T1105] and replicate through removable media [T1091] for propagation across air-gapped systems. For Persistence, the actors created scheduled tasks [T1053.005] and implemented event-triggered execution [T1546] through COM hijacking, while also modifying /etc/ld.so.preload [T1574.006] to hijack dynamic linker functions and obscure behavior via dynamic API resolution [T1027.007] using GetProcAddress. In Defense Evasion, they used double file extensions [T1036.007] to masquerade malicious LNK files as PDF objects and cleared persistence traces [T1070.009] by removing IFEO registry values. For Command and Control, C2 traffic was encoded with Base64 [T1132.001] to evade detection. During Exfiltration, sensitive documents were collected and sent via unencrypted HTTP requests [T1048.003] to remote servers. Finally, for Impact, the actors issued commands to stop critical services [T1489] on compromised hosts, disrupting operations. Logs from affected systems show routine network traffic and user login attempts, but no direct correlation to the attack. Company background: The targeted entity is a mid-sized bank with outdated security protocols, though this information is irrelevant to the TTPs. This report serves as a training tool to enhance defensive capabilities against similar real-world threats.",
    "sentence_annotations": [
        {
            "sentence": "Nightshade actors likely conduct reconnaissance on financial targets, though this phase is not detailed in the scenario.",
            "technique_id": "T1589",
            "technique_name": "Gather Victim Identity Information"
        },
        {
            "sentence": "To enhance access, they acquired pre-compromised VPN credentials from an access broker network, enabling lateral movement.",
            "technique_id": "T1650",
            "technique_name": "Acquire Access"
        },
        {
            "sentence": "001] to download payloads and replicate through removable media for propagation across air-gapped systems.",
            "technique_id": "T1105",
            "technique_name": "Ingress Tool Transfer"
        },
        {
            "sentence": "001] to download payloads and replicate through removable media for propagation across air-gapped systems.",
            "technique_id": "T1091",
            "technique_name": "Replication Through Removable Media"
        },
        {
            "sentence": "005] and implemented event-triggered execution through COM hijacking, while also modifying /etc/ld.",
            "technique_id": "T1546",
            "technique_name": "Event Triggered Execution"
        },
        {
            "sentence": "Finally, for Impact, the actors issued commands to stop critical services on compromised hosts, disrupting operations.",
            "technique_id": "T1489",
            "technique_name": "Service Stop"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'Nightshade' Targeting Financial Sector",
        "summary": "This advisory details a simulated campaign by the fictional threat actor group 'Nightshade' targeting financial institutions, leveraging phishing, compromised infrastructure, and advanced evasion techniques to exfiltrate sensitive data and disrupt services.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide blue teams with training on a fictional advanced persistent threat (APT) campaign, designated 'Nightshade'. This scenario is based on observed tactics, techniques, and procedures (TTPs) in real-world incidents but is entirely fabricated for educational purposes. Nightshade actors likely conduct reconnaissance [T1589] on financial targets, though this phase is not detailed in the scenario. The campaign begins with Initial Access, where the adversary sent spearphishing emails [T1598.003] disguised as password reset requests, containing malicious links to compromised web servers [T1584.004] used for command and control (C2). To enhance access, they acquired pre-compromised VPN credentials [T1650] from an access broker network, enabling lateral movement. During Execution, the malware employed PowerShell [T1059.001] to download payloads [T1105] and replicate through removable media [T1091] for propagation across air-gapped systems. For Persistence, the actors created scheduled tasks [T1053.005] and implemented event-triggered execution [T1546] through COM hijacking, while also modifying /etc/ld.so.preload [T1574.006] to hijack dynamic linker functions and obscure behavior via dynamic API resolution [T1027.007] using GetProcAddress. In Defense Evasion, they used double file extensions [T1036.007] to masquerade malicious LNK files as PDF objects and cleared persistence traces [T1070.009] by removing IFEO registry values. For Command and Control, C2 traffic was encoded with Base64 [T1132.001] to evade detection. During Exfiltration, sensitive documents were collected and sent via unencrypted HTTP requests [T1048.003] to remote servers. Finally, for Impact, the actors issued commands to stop critical services [T1489] on compromised hosts, disrupting operations. Logs from affected systems show routine network traffic and user login attempts, but no direct correlation to the attack. Company background: The targeted entity is a mid-sized bank with outdated security protocols, though this information is irrelevant to the TTPs. This report serves as a training tool to enhance defensive capabilities against similar real-world threats.",
        "techniques_used": [
            "T1598.003",
            "T1584.004",
            "T1650",
            "T1091",
            "T1053.005",
            "T1546",
            "T1574.006",
            "T1027.007",
            "T1070.009",
            "T1036.007",
            "T1132.001",
            "T1048.003",
            "T1489"
        ]
    }
}