# CTI_reports/report_746_20251223_105738_ground_truth.json
{
    "report_id": "report_746",
    "generated_at": "20251223_105738",
    "scenario_target_ids": [
        "T1598",
        "T1593",
        "T1585.001",
        "T1200",
        "T1053.005",
        "T1546.011",
        "T1137.003",
        "T1219.001",
        "T1537"
    ],
    "scenario_flow": [
        "reconnaissance",
        "reconnaissance",
        "resource-development",
        "initial-access",
        "execution",
        "persistence",
        "persistence",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1593",
        "T1585.001",
        "T1598",
        "T1200",
        "T1053.005",
        "T1546.011",
        "T1137.003",
        "T1219.001",
        "T1537"
    ],
    "actual_found_ids": [
        "T1598",
        "T1585.001",
        "T1200",
        "T1137.003",
        "T1546.011",
        "T1593",
        "T1053.005",
        "T1219.001",
        "T1537"
    ],
    "clean_text": "On October 15, 2023, APT-99 initiated a campaign against Acme Corp, a technology firm with 5,000 employees. The actors first conducted reconnaissance by searching open websites and domains, specifically using Facebook to identify recent hires, and then created fake social media accounts with AI-generated profile pictures to gather victim information through tailored spearphishing emails. Acme's security logs showed normal network traffic during this phase, with no alerts triggered. Subsequently, APT-99 gained initial access by planting a USB Rubber Ducky disguised as a standard flash drive in the company breakroom, which an employee inadvertently connected to their workstation. Once executed, the malware established persistence by creating scheduled tasks and applying application shims via a custom.sdb file to disable security features while simultaneously writing a custom Outlook form to the organizational forms library. The adversary then leveraged the compromised development environment to establish a persistent connection through IDE tunneling and used Azure's native backup service to transfer critical file share data to a controlled cloud account. Throughout the incident, Acme's IT team noted unrelated system logs, such as routine software updates and user login attempts from authorized locations, which did not correlate with the attack.",
    "tagged_text": "On October 15, 2023, APT-99 initiated a campaign against Acme Corp, a technology firm with 5,000 employees. The actors first conducted reconnaissance by searching open websites and domains [T1593], specifically using Facebook to identify recent hires, and then created fake social media accounts [T1585.001] with AI-generated profile pictures to gather victim information through tailored spearphishing emails [T1598]. Acme's security logs showed normal network traffic during this phase, with no alerts triggered. Subsequently, APT-99 gained initial access by planting a USB Rubber Ducky disguised as a standard flash drive in the company breakroom [T1200], which an employee inadvertently connected to their workstation. Once executed, the malware established persistence by creating scheduled tasks [T1053.005] and applying application shims via a custom.sdb file [T1546.011] to disable security features while simultaneously writing a custom Outlook form to the organizational forms library [T1137.003]. The adversary then leveraged the compromised development environment to establish a persistent connection through IDE tunneling [T1219.001] and used Azure's native backup service to transfer critical file share data to a controlled cloud account [T1537]. Throughout the incident, Acme's IT team noted unrelated system logs, such as routine software updates and user login attempts from authorized locations, which did not correlate with the attack.",
    "sentence_annotations": [
        {
            "sentence": "The actors first conducted reconnaissance by searching open websites and domains, specifically using Facebook to identify recent hires, and then created fake social media accounts [T1585.",
            "technique_id": "T1593",
            "technique_name": "Search Open Websites/Domains"
        },
        {
            "sentence": "001] with AI-generated profile pictures to gather victim information through tailored spearphishing emails.",
            "technique_id": "T1598",
            "technique_name": "Phishing for Information"
        },
        {
            "sentence": "Subsequently, APT-99 gained initial access by planting a USB Rubber Ducky disguised as a standard flash drive in the company breakroom, which an employee inadvertently connected to their workstation.",
            "technique_id": "T1200",
            "technique_name": "Hardware Additions"
        },
        {
            "sentence": "001] and used Azure's native backup service to transfer critical file share data to a controlled cloud account.",
            "technique_id": "T1537",
            "technique_name": "Transfer Data to Cloud Account"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign Targeting Corporate Networks",
        "summary": "This advisory details a simulated campaign by a fictional threat actor, designated APT-99, targeting a multinational corporation's network. APT-99 employed sophisticated techniques including social engineering, hardware-based initial access, and persistence mechanisms to exfiltrate sensitive data. The report is intended for blue team training to enhance detection and response capabilities.",
        "full_report_text": "On October 15, 2023, APT-99 initiated a campaign against Acme Corp, a technology firm with 5,000 employees. The actors first conducted reconnaissance by searching open websites and domains [T1593], specifically using Facebook to identify recent hires, and then created fake social media accounts [T1585.001] with AI-generated profile pictures to gather victim information through tailored spearphishing emails [T1598]. Acme's security logs showed normal network traffic during this phase, with no alerts triggered. Subsequently, APT-99 gained initial access by planting a USB Rubber Ducky disguised as a standard flash drive in the company breakroom [T1200], which an employee inadvertently connected to their workstation. Once executed, the malware established persistence by creating scheduled tasks [T1053.005] and applying application shims via a custom.sdb file [T1546.011] to disable security features while simultaneously writing a custom Outlook form to the organizational forms library [T1137.003]. The adversary then leveraged the compromised development environment to establish a persistent connection through IDE tunneling [T1219.001] and used Azure's native backup service to transfer critical file share data to a controlled cloud account [T1537]. Throughout the incident, Acme's IT team noted unrelated system logs, such as routine software updates and user login attempts from authorized locations, which did not correlate with the attack.",
        "techniques_used": [
            "T1593",
            "T1585.001",
            "T1598",
            "T1200",
            "T1053.005",
            "T1546.011",
            "T1137.003",
            "T1219.001",
            "T1537"
        ]
    }
}