# CTI_reports/report_386_20251222_152519_ground_truth.json
{
    "report_id": "report_386",
    "generated_at": "20251222_152519",
    "scenario_target_ids": [
        "T1596.003",
        "T1566.002",
        "T1053",
        "T1546.012",
        "T1574.013",
        "T1547.010",
        "T1574.014",
        "T1205.001",
        "T1561"
    ],
    "scenario_flow": [
        "reconnaissance",
        "initial-access",
        "execution",
        "persistence",
        "persistence",
        "privilege-escalation",
        "privilege-escalation",
        "command-and-control",
        "impact"
    ],
    "llm_declared_ids": [
        "T1596.003",
        "T1566.002",
        "T1053",
        "T1546.012",
        "T1574.013",
        "T1547.010",
        "T1574.014",
        "T1205.001",
        "T1561"
    ],
    "actual_found_ids": [
        "T1574.014",
        "T1566.002",
        "T1547.010",
        "T1596.003",
        "T1205.001",
        "T1546.012",
        "T1574.013",
        "T1053",
        "T1561"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight a fictional APT campaign, 'SilentWipe', designed for blue team training. This scenario simulates an attack on a mid-sized financial services firm, 'SecureBank Inc.', based in New York, with operations across North America. The firm's network logs showed unusual DNS queries from an internal workstation at 10:00 UTC on April 15, 2024, but no immediate alerts were triggered.\n\nInitial Access: The actors likely conducted reconnaissance by querying Certificate Transparency logs to identify subdomains associated with SecureBank, such as 'dev.securebank.com'. Following this, they sent spearphishing emails to targeted employees, masquerading as a vendor update with a malicious link that, when clicked, downloaded a payload disguised as a PDF document. The email was sent from a spoofed domain, 'support@secure-vendor.com', and bypassed initial spam filters due to its personalized content.\n\nExecution: Upon execution, the malware employed PowerShell to schedule a task that ran 'data_exfil.exe' every 30 minutes, while simultaneously modifying registry keys in HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\SilentProcessExit\\services.exe to trigger malicious code when critical services stopped. This was combined with injecting code into a legitimate browser process using WriteProcessMemory to evade detection by security software.\n\nPersistence: To maintain access, the adversary wrote a malicious DLL named 'monsvc.dll' to the Windows system folder for port monitor hijacking and modified the application's .config file to redirect assembly loading to malicious code via AppDomainManager. Additionally, they used raw sockets for port knocking to communicate with the C2 server without standard protocols, ensuring stealthy command and control channels.\n\nImpact: In the final stage, the actors deployed a destructive script that systematically erased backup files stored on network shares, causing data loss and disrupting recovery efforts. SecureBank's incident response team noted that the script targeted '.bak' and '.vhd' files across multiple servers, with logs indicating deletion activities starting at 14:30 UTC. No actual financial data was exfiltrated in this training scenario, but the simulation highlighted vulnerabilities in backup management.\n\nThis report is based on fictional data for educational purposes. Blue teams should review logs for anomalies in scheduled tasks, registry modifications, and unusual network traffic patterns to enhance detection capabilities.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight a fictional APT campaign, 'SilentWipe', designed for blue team training. This scenario simulates an attack on a mid-sized financial services firm, 'SecureBank Inc.', based in New York, with operations across North America. The firm's network logs showed unusual DNS queries from an internal workstation at 10:00 UTC on April 15, 2024, but no immediate alerts were triggered.\n\nInitial Access: The actors likely conducted reconnaissance by querying Certificate Transparency logs [T1596.003] to identify subdomains associated with SecureBank, such as 'dev.securebank.com'. Following this, they sent spearphishing emails [T1566.002] to targeted employees, masquerading as a vendor update with a malicious link that, when clicked, downloaded a payload disguised as a PDF document. The email was sent from a spoofed domain, 'support@secure-vendor.com', and bypassed initial spam filters due to its personalized content.\n\nExecution: Upon execution, the malware employed PowerShell to schedule a task [T1053] that ran 'data_exfil.exe' every 30 minutes, while simultaneously modifying registry keys in HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\SilentProcessExit\\services.exe [T1546.012] to trigger malicious code when critical services stopped. This was combined with injecting code into a legitimate browser process using WriteProcessMemory [T1574.013] to evade detection by security software.\n\nPersistence: To maintain access, the adversary wrote a malicious DLL named 'monsvc.dll' to the Windows system folder for port monitor hijacking [T1547.010] and modified the application's .config file to redirect assembly loading to malicious code via AppDomainManager [T1574.014]. Additionally, they used raw sockets for port knocking [T1205.001] to communicate with the C2 server without standard protocols, ensuring stealthy command and control channels.\n\nImpact: In the final stage, the actors deployed a destructive script that systematically erased backup files stored on network shares [T1561], causing data loss and disrupting recovery efforts. SecureBank's incident response team noted that the script targeted '.bak' and '.vhd' files across multiple servers, with logs indicating deletion activities starting at 14:30 UTC. No actual financial data was exfiltrated in this training scenario, but the simulation highlighted vulnerabilities in backup management.\n\nThis report is based on fictional data for educational purposes. Blue teams should review logs for anomalies in scheduled tasks, registry modifications, and unusual network traffic patterns to enhance detection capabilities.",
    "sentence_annotations": [
        {
            "sentence": "Execution: Upon execution, the malware employed PowerShell to schedule a task that ran 'data_exfil.",
            "technique_id": "T1053",
            "technique_name": "Scheduled Task/Job"
        },
        {
            "sentence": "Impact: In the final stage, the actors deployed a destructive script that systematically erased backup files stored on network shares, causing data loss and disrupting recovery efforts.",
            "technique_id": "T1561",
            "technique_name": "Disk Wipe"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'SilentWipe' Targeting Financial Sector",
        "summary": "This advisory details a simulated campaign by a fictional APT group, 'SilentWipe', targeting financial institutions. The campaign involves reconnaissance via digital certificates, spearphishing for initial access, and sophisticated persistence mechanisms leading to data exfiltration and destructive disk wiping. Blue teams should use this report for training to identify and mitigate similar tactics.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to highlight a fictional APT campaign, 'SilentWipe', designed for blue team training. This scenario simulates an attack on a mid-sized financial services firm, 'SecureBank Inc.', based in New York, with operations across North America. The firm's network logs showed unusual DNS queries from an internal workstation at 10:00 UTC on April 15, 2024, but no immediate alerts were triggered.\n\nInitial Access: The actors likely conducted reconnaissance by querying Certificate Transparency logs [T1596.003] to identify subdomains associated with SecureBank, such as 'dev.securebank.com'. Following this, they sent spearphishing emails [T1566.002] to targeted employees, masquerading as a vendor update with a malicious link that, when clicked, downloaded a payload disguised as a PDF document. The email was sent from a spoofed domain, 'support@secure-vendor.com', and bypassed initial spam filters due to its personalized content.\n\nExecution: Upon execution, the malware employed PowerShell to schedule a task [T1053] that ran 'data_exfil.exe' every 30 minutes, while simultaneously modifying registry keys in HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\SilentProcessExit\\services.exe [T1546.012] to trigger malicious code when critical services stopped. This was combined with injecting code into a legitimate browser process using WriteProcessMemory [T1574.013] to evade detection by security software.\n\nPersistence: To maintain access, the adversary wrote a malicious DLL named 'monsvc.dll' to the Windows system folder for port monitor hijacking [T1547.010] and modified the application's .config file to redirect assembly loading to malicious code via AppDomainManager [T1574.014]. Additionally, they used raw sockets for port knocking [T1205.001] to communicate with the C2 server without standard protocols, ensuring stealthy command and control channels.\n\nImpact: In the final stage, the actors deployed a destructive script that systematically erased backup files stored on network shares [T1561], causing data loss and disrupting recovery efforts. SecureBank's incident response team noted that the script targeted '.bak' and '.vhd' files across multiple servers, with logs indicating deletion activities starting at 14:30 UTC. No actual financial data was exfiltrated in this training scenario, but the simulation highlighted vulnerabilities in backup management.\n\nThis report is based on fictional data for educational purposes. Blue teams should review logs for anomalies in scheduled tasks, registry modifications, and unusual network traffic patterns to enhance detection capabilities.",
        "techniques_used": [
            "T1596.003",
            "T1566.002",
            "T1053",
            "T1546.012",
            "T1574.013",
            "T1547.010",
            "T1574.014",
            "T1205.001",
            "T1561"
        ]
    }
}