# generated_reports/report_230_20251222_141954_ground_truth.json
{
    "report_id": "report_230",
    "generated_at": "20251222_141954",
    "scenario_target_ids": [
        "T1594",
        "T1190",
        "T1129",
        "T1059.009",
        "T1053.005",
        "T1552.005",
        "T1092",
        "T1102.001"
    ],
    "scenario_flow": [
        "reconnaissance",
        "initial-access",
        "execution",
        "execution",
        "persistence",
        "credential-access",
        "command-and-control",
        "command-and-control"
    ],
    "llm_declared_ids": [
        "T1594",
        "T1190",
        "T1129",
        "T1059.009",
        "T1053.005",
        "T1552.005",
        "T1092",
        "T1102.001"
    ],
    "actual_found_ids": [
        "T1129",
        "T1102.001",
        "T1552.005",
        "T1059.009",
        "T1092",
        "T1190",
        "T1594",
        "T1053.005"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide blue teams with training on a fictional APT campaign. 'ShadowGate' actors likely conduct reconnaissance by searching victim-owned websites to gather intelligence on internal network structures and exposed services. Following this, they exploit a vulnerable FortiGate device to gain initial access, deploying malware that utilizes LoadLibraryEx to load a rogue DLL for establishing persistence while simultaneously creating scheduled tasks to maintain long-term access. Once inside, the actors employ Azure PowerShell to remotely execute commands like 'Get-AzVM' for enumerating virtual machines and query the AWS instance metadata service to steal cloud credentials. To facilitate command and control across air-gapped networks, they write configuration files to USB drives and use Pastebin-hosted webpages as dead drop resolvers to update C2 domains. System logs from affected organizations show increased network traffic to internal web servers, but no unauthorized access was detected in this simulated scenario. The company, a mid-sized tech firm with hybrid AWS and Azure deployments, reported no actual data loss, as this is a training exercise.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide blue teams with training on a fictional APT campaign. 'ShadowGate' actors likely conduct reconnaissance by searching victim-owned websites [T1594] to gather intelligence on internal network structures and exposed services. Following this, they exploit a vulnerable FortiGate device [T1190] to gain initial access, deploying malware that utilizes LoadLibraryEx to load a rogue DLL [T1129] for establishing persistence while simultaneously creating scheduled tasks [T1053.005] to maintain long-term access. Once inside, the actors employ Azure PowerShell [T1059.009] to remotely execute commands like 'Get-AzVM' for enumerating virtual machines and query the AWS instance metadata service [T1552.005] to steal cloud credentials. To facilitate command and control across air-gapped networks, they write configuration files to USB drives [T1092] and use Pastebin-hosted webpages as dead drop resolvers [T1102.001] to update C2 domains. System logs from affected organizations show increased network traffic to internal web servers, but no unauthorized access was detected in this simulated scenario. The company, a mid-sized tech firm with hybrid AWS and Azure deployments, reported no actual data loss, as this is a training exercise.",
    "sentence_annotations": [
        {
            "sentence": "'ShadowGate' actors likely conduct reconnaissance by searching victim-owned websites to gather intelligence on internal network structures and exposed services.",
            "technique_id": "T1594",
            "technique_name": "Search Victim-Owned Websites"
        },
        {
            "sentence": "Following this, they exploit a vulnerable FortiGate device to gain initial access, deploying malware that utilizes LoadLibraryEx to load a rogue DLL for establishing persistence while simultaneously creating scheduled tasks [T1053.",
            "technique_id": "T1190",
            "technique_name": "Exploit Public-Facing Application"
        },
        {
            "sentence": "Following this, they exploit a vulnerable FortiGate device to gain initial access, deploying malware that utilizes LoadLibraryEx to load a rogue DLL for establishing persistence while simultaneously creating scheduled tasks [T1053.",
            "technique_id": "T1129",
            "technique_name": "Shared Modules"
        },
        {
            "sentence": "To facilitate command and control across air-gapped networks, they write configuration files to USB drives and use Pastebin-hosted webpages as dead drop resolvers [T1102.",
            "technique_id": "T1092",
            "technique_name": "Communication Through Removable Media"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'ShadowGate' Targeting Cloud and On-Premises Infrastructure",
        "summary": "This advisory details a simulated campaign by a fictional advanced persistent threat (APT) group, designated 'ShadowGate', targeting organizations with hybrid cloud environments. The campaign involves reconnaissance, exploitation of public-facing applications, use of shared modules for persistence, cloud API abuse, scheduled tasks, metadata API queries, removable media communication, and dead drop resolvers for command and control. This report is for educational purposes to train blue teams on identifying and mitigating such tactics.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide blue teams with training on a fictional APT campaign. 'ShadowGate' actors likely conduct reconnaissance by searching victim-owned websites [T1594] to gather intelligence on internal network structures and exposed services. Following this, they exploit a vulnerable FortiGate device [T1190] to gain initial access, deploying malware that utilizes LoadLibraryEx to load a rogue DLL [T1129] for establishing persistence while simultaneously creating scheduled tasks [T1053.005] to maintain long-term access. Once inside, the actors employ Azure PowerShell [T1059.009] to remotely execute commands like 'Get-AzVM' for enumerating virtual machines and query the AWS instance metadata service [T1552.005] to steal cloud credentials. To facilitate command and control across air-gapped networks, they write configuration files to USB drives [T1092] and use Pastebin-hosted webpages as dead drop resolvers [T1102.001] to update C2 domains. System logs from affected organizations show increased network traffic to internal web servers, but no unauthorized access was detected in this simulated scenario. The company, a mid-sized tech firm with hybrid AWS and Azure deployments, reported no actual data loss, as this is a training exercise.",
        "techniques_used": [
            "T1594",
            "T1190",
            "T1129",
            "T1059.009",
            "T1053.005",
            "T1552.005",
            "T1092",
            "T1102.001"
        ]
    }
}