# generated_reports/report_541_20251222_163256_ground_truth.json
{
    "report_id": "report_541",
    "generated_at": "20251222_163256",
    "scenario_target_ids": [
        "T1589.002",
        "T1566",
        "T1053.006",
        "T1053.005",
        "T1053.007",
        "T1218.008",
        "T1127.003",
        "T1021.002",
        "T1119",
        "T1039",
        "T1114.001",
        "T1095",
        "T1048"
    ],
    "scenario_flow": [
        "reconnaissance",
        "initial-access",
        "execution",
        "execution",
        "persistence",
        "defense-evasion",
        "defense-evasion",
        "lateral-movement",
        "collection",
        "collection",
        "collection",
        "command-and-control",
        "exfiltration"
    ],
    "llm_declared_ids": [
        "T1589.002",
        "T1566",
        "T1053.006",
        "T1053.005",
        "T1053.007",
        "T1218.008",
        "T1127.003",
        "T1021.002",
        "T1119",
        "T1039",
        "T1114.001",
        "T1095",
        "T1048"
    ],
    "actual_found_ids": [
        "T1039",
        "T1095",
        "T1053.006",
        "T1114.001",
        "T1127.003",
        "T1566",
        "T1053.007",
        "T1589.002",
        "T1119",
        "T1021.002",
        "T1218.008",
        "T1053.005",
        "T1048"
    ],
    "clean_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide insights into a fictional APT campaign for educational purposes. Actors likely conduct reconnaissance by querying APIs to identify corporate email addresses, then send phishing emails with malicious macro-laden invoice attachments to gain initial access. Upon execution, the adversary employs systemd timers and scheduled tasks for persistence, while using container orchestration jobs to deploy malicious containers that establish command and control (C2) via TCP. To evade detection, they execute a malicious DLL via odbcconf.exe and use JamPlus.exe to run scripts that exfiltrate data. For lateral movement, the actors access SMB/Windows admin shares to copy payloads and schedule execution, while enumerating network shares with tools like 'net view' to identify accessible file servers. Data collection involves automated tools to gather system and network information and commands to collect emails from Outlook. Finally, exfiltration occurs over alternative protocols, such as NetBIOS over TCP, to bypass network monitoring. Logs from the compromised systems show routine administrative activities unrelated to the attack, such as system updates and user logins, which were not tagged as malicious. The fictional company, 'FinSecure Bank', reported no actual financial loss, but this scenario highlights common attack vectors in the financial sector.",
    "tagged_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide insights into a fictional APT campaign for educational purposes. Actors likely conduct reconnaissance by querying APIs to identify corporate email addresses [T1589.002], then send phishing emails with malicious macro-laden invoice attachments [T1566] to gain initial access. Upon execution, the adversary employs systemd timers [T1053.006] and scheduled tasks [T1053.005] for persistence, while using container orchestration jobs [T1053.007] to deploy malicious containers that establish command and control (C2) via TCP [T1095]. To evade detection, they execute a malicious DLL via odbcconf.exe [T1218.008] and use JamPlus.exe to run scripts that exfiltrate data [T1127.003]. For lateral movement, the actors access SMB/Windows admin shares [T1021.002] to copy payloads and schedule execution, while enumerating network shares with tools like 'net view' [T1039] to identify accessible file servers. Data collection involves automated tools to gather system and network information [T1119] and commands to collect emails from Outlook [T1114.001]. Finally, exfiltration occurs over alternative protocols, such as NetBIOS over TCP [T1048], to bypass network monitoring. Logs from the compromised systems show routine administrative activities unrelated to the attack, such as system updates and user logins, which were not tagged as malicious. The fictional company, 'FinSecure Bank', reported no actual financial loss, but this scenario highlights common attack vectors in the financial sector.",
    "sentence_annotations": [
        {
            "sentence": "002], then send phishing emails with malicious macro-laden invoice attachments to gain initial access.",
            "technique_id": "T1566",
            "technique_name": "Phishing"
        },
        {
            "sentence": "007] to deploy malicious containers that establish command and control (C2) via TCP.",
            "technique_id": "T1095",
            "technique_name": "Non-Application Layer Protocol"
        },
        {
            "sentence": "002] to copy payloads and schedule execution, while enumerating network shares with tools like 'net view' to identify accessible file servers.",
            "technique_id": "T1039",
            "technique_name": "Data from Network Shared Drive"
        },
        {
            "sentence": "Data collection involves automated tools to gather system and network information and commands to collect emails from Outlook [T1114.",
            "technique_id": "T1119",
            "technique_name": "Automated Collection"
        },
        {
            "sentence": "Finally, exfiltration occurs over alternative protocols, such as NetBIOS over TCP, to bypass network monitoring.",
            "technique_id": "T1048",
            "technique_name": "Exfiltration Over Alternative Protocol"
        }
    ],
    "llm_raw_response": {
        "title": "Joint Cybersecurity Advisory: Fictional APT Campaign 'ShadowHarvest' Targeting Financial Sector",
        "summary": "This advisory details a simulated campaign by a fictional advanced persistent threat (APT) group, designated as 'ShadowHarvest', targeting financial institutions. The campaign involves initial reconnaissance, phishing, persistence mechanisms, lateral movement, data collection, and exfiltration using various techniques. Blue teams can use this report for training to identify and mitigate similar attack patterns.",
        "full_report_text": "The Federal Bureau of Investigation (FBI), CISA, and NSA are releasing this joint Cybersecurity Advisory to provide insights into a fictional APT campaign for educational purposes. Actors likely conduct reconnaissance by querying APIs to identify corporate email addresses [T1589.002], then send phishing emails with malicious macro-laden invoice attachments [T1566] to gain initial access. Upon execution, the adversary employs systemd timers [T1053.006] and scheduled tasks [T1053.005] for persistence, while using container orchestration jobs [T1053.007] to deploy malicious containers that establish command and control (C2) via TCP [T1095]. To evade detection, they execute a malicious DLL via odbcconf.exe [T1218.008] and use JamPlus.exe to run scripts that exfiltrate data [T1127.003]. For lateral movement, the actors access SMB/Windows admin shares [T1021.002] to copy payloads and schedule execution, while enumerating network shares with tools like 'net view' [T1039] to identify accessible file servers. Data collection involves automated tools to gather system and network information [T1119] and commands to collect emails from Outlook [T1114.001]. Finally, exfiltration occurs over alternative protocols, such as NetBIOS over TCP [T1048], to bypass network monitoring. Logs from the compromised systems show routine administrative activities unrelated to the attack, such as system updates and user logins, which were not tagged as malicious. The fictional company, 'FinSecure Bank', reported no actual financial loss, but this scenario highlights common attack vectors in the financial sector.",
        "techniques_used": [
            "T1589.002",
            "T1566",
            "T1053.006",
            "T1053.005",
            "T1053.007",
            "T1218.008",
            "T1127.003",
            "T1021.002",
            "T1119",
            "T1039",
            "T1114.001",
            "T1095",
            "T1048"
        ]
    }
}